# crontab
* * * * * . /environment; /usr/bin/python3 /data.py >> /var/log/cron.log
